# config/etc/rc.local
sudo python /home/pi/lynniphone/src/init/init.py &
exit 0